# files/fedora.manifest
phc2sys.sysconfig /etc/sysconfig/phc2sys
ptp4l-gm.service /etc/systemd/system/ptp4l-gm.service
ptp4l.conf /etc/ptp4l.conf
refclock-shm.conf /etc/chrony.d/refclock-shm.conf
ts2phc.conf /etc/ts2phc.conf
ts2phc.service /etc/systemd/system/ts2phc.service
ts2phc.sysconfig /etc/sysconfig/ts2phc
